# .config/systemd/user/checkupdates.timer
[Unit]
Description=Check for pacman and AUR updates hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
