# .config/systemd/user/checkupdates.timer
[Unit]
Description=Check for pacman and AUR updates daily

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true
