# .config/systemd/user/journal.service
[Unit]
Description=Update journal indicator in status bar

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -RTMIN+13 i3blocks
