# .config/systemd/user/mbsync.service
[Unit]
Description=Sync emails

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
