# .config/systemd/user/mbsync.service
[Unit]
Description=Sync emails

[Service]
Type=oneshot
ExecStart=/home/llama/bin/ms
