# .config/systemd/user/mbsync.service
[Unit]
Description=Sync emails

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'killall -q mbsync; ( mbsync -Va; notmuch new ) & wait'
