# .config/systemd/user/mbsync.service
[Unit]
Description=Sync emails

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'mbsync -Va; notmuch new'
