# .config/systemd/user/mbsync.timer
[Unit]
Description=Sync emails minutely

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
