# .config/systemd/user/wallpaper.timer
[Unit]
Description=Randomize the wallpaper hourly

[Timer]
OnCalendar=hourly
Persistent=true
